# service/repostbot.service
[Unit]
Wants=repostbot.timer

[Service]
Type=oneshot
EnvironmentFile=
ExecStart=/usr/bin/env repostbot

[Install]
WantedBy=timers.target
